# consul/etc/services.d/crond/run
#!/usr/bin/with-contenv sh
set -e

exec crond -b -L /dev/null