# consul/etc/services.d/nginx/finish
#!/usr/bin/with-contenv sh
set -e

s6-svscanctl -t /var/run/s6/services