# consul/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh
set -e

NGINX_CONF=${NGINX_CONF:-/etc/nginx/nginx.conf}
nginx -c $NGINX_CONF -t && exec nginx -c $NGINX_CONF -g "daemon off;"