# consul/etc/services.d/consul-template/run
#!/usr/bin/with-contenv sh
set -e

if [ -n "$CONSUL_URL" -a -d /etc/consul-template/conf ]; then
    exec consul-template -consul ${CONSUL_URL} -config "/etc/consul-template/conf"
fi