# consul/etc/crontabs/root
# do daily/weekly/monthly maintenance
# min   hour    day     month   weekday command
*/1    *       *       *       *       run-parts /etc/periodic/1min
